# src/CronBundle/Infrastucture/Repository/CronFiles/crontab
*/10 * * * * /scripts/monitor.sh
0 2 * * * /bin/sh backup.sh
@daily /scripts/script.sh
3-59/15 6-12 */15 1 2-5 ls
0 4,17 * * sun,mon /scripts/script.sh
